# minecraft/roles/setup/files/minecraft.service
[Unit]
Description=Minecraft Server
After=network.target

[Service]
User=darshan
Group=darshan
WorkingDirectory=/etc/minecraft

ExecStart=/usr/bin/java -Xms2G -Xmx3G -jar /etc/minecraft/server.jar nogui
Restart=on-failure
RestartSec=5
Nice=5

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
ReadWritePaths=/etc/minecraft

[Install]
WantedBy=multi-user.target